# services/ntfy.service
[Unit]
Description=Signal ntfy service

[Service]
User=ryandterri
ExecStart=/usr/local/bin/signal_receive

[Install]
WantedBy=default.target
