# services/ntfy.timer
[Unit]
Description=Signal ntfy service

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target
