# services/nzxt.service
[Unit]
Description=NZXT gif loader

[Service]
ExecStart=/usr/local/bin/nzxt_setup

[Install]
WantedBy=default.target
