# services/open_rgb.timer
[Unit]
Description=Set rgb color

[Timer]
OnBootSec=20

[Install]
WantedBy=timers.target
