# services/open_rgb.service
[Unit]
Description=Set RGB colors

[Service]
ExecStart=/usr/local/bin/set_rgb

[Install]
WantedBy=default.target
